# systemd/emis-node.service
[Unit]
Description=ÉMIS Community Sensor Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=emis
WorkingDirectory=/home/emis/emis/code
ExecStart=/usr/bin/python3 /home/emis/emis/code/collect_data.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

